# code/rootfs/rootfs/etc/init.d/rc.local
#!/bin/sh
# insert kernel modules
# modprobe xxx

mknod /dev/random c 1 8
mknod /dev/urandom c 1 9