# code/rootfs/rootfs/etc/init.d/rc.local
#!/bin/sh
# insert kernel modules
# modprobe xxx

